# com.heartflow.auto-upgrade.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.heartflow.auto-upgrade</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/apple/mark-heartflow-skill/scripts/auto-upgrade.js</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>MinuteInterval</key>
            <integer>30</integer>
        </dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/apple/mark-heartflow-skill/logs/auto-upgrade.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/apple/mark-heartflow-skill/logs/auto-upgrade.error.log</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>