# com.heartflow.paper-upgrade.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.heartflow.paper-upgrade</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-l</string>
        <string>-c</string>
        <string>export PATH="/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin"; /Users/apple/.hermes/skills/ai/mark-heartflow-skill/run-paper-upgrade.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/apple/.hermes/skills/ai/mark-heartflow-skill/logs/launchd.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/apple/.hermes/skills/ai/mark-heartflow-skill/logs/launchd.stderr.log</string>
</dict>
</plist>
