# com.heartflow.paper-upgrade.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.heartflow.paper-upgrade</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/apple/.hermes/skills/ai/mark-heartflow-skill/cron/paper-upgrade-engine.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>3600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/apple/.hermes/skills/ai/mark-heartflow-skill/logs/paper-upgrade.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/apple/.hermes/skills/ai/mark-heartflow-skill/logs/paper-upgrade.stderr.log</string>
</dict>
</plist>
