# evolution/com.heartflow.paper-analyzer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.heartflow.paper-analyzer</string>

    <!-- 每 12 小时：凌晨和中午 -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>0</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>12</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
    </array>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/apple/.hermes/skills/ai/heartflow/evolution/run_analyze.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/apple/.hermes/skills/ai/heartflow/evolution</string>

    <key>StandardOutPath</key>
    <string>/Users/apple/.hermes/skills/ai/heartflow/evolution_logs/analyze_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/apple/.hermes/skills/ai/heartflow/evolution_logs/analyze_stderr.log</string>

    <key>LowPriorityBackgroundIO</key>
    <true/>
</dict>
</plist>